# templates/supervisord.conf.j2
; This file was generated by Ansible for {{ ansible_fqdn }}
; Do NOT modify this file by hand!

; supervisor config file

[inet_http_server]
port = 127.0.0.1:9001

[supervisord]
logfile={{supervisor_logdir}}/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile={{supervisor_pid}}; (supervisord pidfile;default supervisord.pid)
childlogdir={{supervisor_logdir}}            ; ('AUTO' child log dir, default $TEMP)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

; programs
{% for task in supervisor_tasks %}
[program:{{task.name}}]
{% for option, value in task.items() %}
{%if option != 'name' %}
{{option}} = {{value}}
{% endif %}
{% endfor %}
{% endfor %}

; The [include] section can just contain the "files" setting.  This
; setting can list multiple files (separated by whitespace or
; newlines).  It can also contain wildcards.  The filenames are
; interpreted as relative to this file.  Included files *cannot*
; include files themselves.

[include]
files = {{supervisor_incdir}}/*.conf
